# yocto/sources/meta-gabriel/recipes-core/systemd/files/block-device-detach@.service
[Unit]
Description=Remove block device

[Service]
Type=oneshot
ExecStart=/opt/gabriel/bin/block-device-detach.sh %i &

[Install]
WantedBy=multi-user.target
