# yocto/sources/meta-gabriel/recipes-core/systemd/files/block-device-attach@.service
[Unit]
Description=Add block device

[Service]
Type=oneshot
ExecStart=/opt/gabriel/bin/block-device-attach.sh %i &

[Install]
WantedBy=multi-user.target
