# yocto/sources/meta-gabriel/recipes-core/systemd/files/xdisk.service
[Unit]
Description=Pingu's script to expand disk (rootfs)

[Service]
ExecStart=/usr/bin/xdisk

[Install]
WantedBy=sysinit.target